# deploy/systemd/web.service
[Unit]
Description=Dance Studio Web
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dance
Group=dance
WorkingDirectory=/opt/dance_studio
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/opt/dance_studio/src"
EnvironmentFile=-/opt/dance_studio/.env
ExecStart=/opt/dance_studio/.venv/bin/gunicorn dance_studio.web.wsgi:app --bind 127.0.0.1:3000 --workers 2 --threads 4 --timeout 30 --graceful-timeout 30
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
